# build/gh-wh-handler.service.in
[Unit]
Description=Runs github webhook handler
After=network.target

[Service]
ExecStart=@SERVICE_EXECUTABLE@ @SERVICE_CONFIG@ @SERVICE_LOGS@
Restart=always
Type=simple

[Install]
WantedBy=default.target
